# etc/s6-overlay/s6-rc.d/litestream/run
#!/command/with-contenv bash

exec litestream replicate
