# etc/services.d/litestream/run
#!/usr/bin/with-contenv bash

exec litestream replicate
