# etc/s6-overlay/s6-rc.d/myapp/run
#!/command/with-contenv bash
exec myapp -dsn /data/db
