# etc/services.d/myapp/finish
#!/usr/bin/execlineb -S0

# Stops all services when "myapp" stops.
s6-svscanctl -t /var/run/s6/services
