# etc/services.d/myapp/run
#!/usr/bin/with-contenv bash
/usr/bin/myapp -dsn /data/db
