# etc/services.d/myapp/run
#!/usr/bin/with-contenv bash
myapp -dsn /data/db
