# etc/s6-overlay/s6-rc.d/litestream-restore/up
/etc/s6-overlay/scripts/litestream-restore